# docker/crontab
2,7,12,17,22,27,32,37,42,47,52,57 9,10,11,12,13,14,15,16 * * 2,3,4,5 root /script.sh > /proc/1/fd/1 2>/proc/1/fd/2
# empty line